# sma-em/rootfs/etc/services.d/sma-em/run
#!/command/execlineb -P

uv run python -m ha_addon_sma
