# sma-em/rootfs/etc/s6-overlay/s6-rc.d/sma-em/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting the SMA-EM add-on"
uv run python -m ha_addon_sma
