# sma-em-dev/rootfs/etc/services.d/sma-em/run
#!/command/execlineb -P

python3 /usr/src/app/run.py
